# dist/debian/mariner3d.service
[Unit]
Description=mapache3d Printer Controller

[Service]
User=mapache
ExecStart=/usr/bin/mapache
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
